# files/launch/setup-environment.conf.in
description "Setup environment for all other boot jobs"

start on startup

script
    echo "no op"
end script
