# files/launch/udev.conf.in
description "device node and kernel event manager"

start on stopped mountall

exec /etc/init.d/udev start
